# crontab
# Fetch from digg to run every 5 mins, but minDelay should be 30 mins. 6 attempts if it fails.
*/5  * * * * /Users/Bob/Sites/imagepush2/app/console imagepush:fetch-from-digg
